# install/crontab/crontab.template
#CUCKOOCRONS
0 * * * * CUCKOODIR/utils/crontab/cuckoo/clean_mem.sh
#ENDCUCKOOCRONS